# launchd/com.USER.claude-respawn.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.__USER__.claude-respawn</string>
	<key>ProgramArguments</key>
	<array>
		<string>__HOME__/bin/claude-respawn</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-respawn.err</string>
</dict>
</plist>
