# launchd/com.USER.claude-sessions-refresh.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.__USER__.claude-sessions-refresh</string>
	<key>ProgramArguments</key>
	<array>
		<string>__HOME__/bin/claude-sessions-refresh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-sessions-refresh.err</string>
</dict>
</plist>
